# deploy/openinbound-backup.timer
[Unit]
Description=Daily backup of OpenInboundEmail data

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
